# dockers/pdns/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:pdns]
command=/usr/sbin/pdns_server --config-dir=/etc/ --master --daemon=no

[program:request]
command=/usr/bin/ruby /rcmd/request.rb
